# deployment/upstart/uwsgi.conf
description "uwsgi starter"
start on runlevel [2345]
stop on runlevel [06]

respawn

exec uwsgi --ini /var/www/stream/osp_uwsgi.ini